# truss/templates/docker_server/supervisord.conf.jinja
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:model-server]
command={{start_command}}
startsecs=30
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[eventlistener:quit_on_failure]
events=PROCESS_STATE_FATAL
command=sh -c 'echo "READY"; read line; kill -15 1; echo "RESULT 2";'

